# avsecure.service
[Unit]
Description=AVSecure WebRTC signaling server
After=network.target

[Service]
Type=simple
User=avsecure
Group=avsecure
WorkingDirectory=/opt/avsecure
ExecStart=/opt/avsecure/go-react-webrtc
Restart=on-failure
RestartSec=5

# Systemd hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/avsecure
PrivateTmp=true
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
